# deploy/canary/farol/ai-token-soak-audit.timer
[Unit]
Description=Daily cumulative ai-token live-soak audit

[Timer]
OnCalendar=*-*-* 04:40:00 UTC
Persistent=true
Unit=ai-token-soak-audit.service

[Install]
WantedBy=timers.target
