# deploy/canary/farol/ai-token-canary.service
[Unit]
Description=Two-hour ai-token leader lifecycle canary
After=network-online.target ai-token-canary-vault.service
Wants=network-online.target ai-token-canary-vault.service
OnFailure=ai-token-canary-alert@%n.service

[Service]
Type=oneshot
ExecStart=/home/kelvin/.ai-token-canary/release/current/tools/run-live-canary --config /home/kelvin/.ai-token-canary/canary.json --live
