# deploy/canary/farol/ai-token-canary.timer
[Unit]
Description=Schedule daily ai-token leader lifecycle canary

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
